# brisa-api.service
[Unit]
Description=Brisa Path Search API Server
After=network.target

[Service]
Type=simple
User=bitiz
WorkingDirectory=/home/bitiz/Desktop/brisa_endpoint
Environment="PATH=/home/bitiz/Desktop/brisa_endpoint/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="FLASK_DEBUG=False"
ExecStart=/bin/bash /home/bitiz/Desktop/brisa_endpoint/start-api.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
